# scripts/vw-mysql-to-sqlite.timer
[Unit]
Description=Hourly Vaultwarden warm-standby sync

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Unit=vw-mysql-to-sqlite.service

[Install]
WantedBy=timers.target
